# gsheets_test_tg/bot.conf
[program:akfa]
directory=/home/ubuntu/pygsheets/gsheets_test_tg
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu", USER="ubuntu"

stderr_logfile=/home/ubuntu/pygsheets/logfile_err_with_pygsheets.log
stdout_logfile=/home/ubuntu/pygsheets/logfile_pygsheets.log
